# docker/php/crontab
0 * * * * php /var/www/bin/console app:wind-forecast --steps=24 --env=prod --no-debug --no-interaction
